# email-cron/config/crontab
# Run the email script every week
0 0 * * 0 /app/run_email.sh
