# email-cron/config/crontab
# Run the email script every minute
* * * * * /app/run_email.sh
