# cron_email/config/crontab
# Run the email script every hour
0 * * * * cd /app && /opt/conda/envs/email_env/bin/python /app/send_email.py >> /var/log/cron.log 2>&1
